# setup_native/source/opensolaris/bundledextensions/ooo_bundled_extensions.xml
<?xml version="1.0"?>

<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type="manifest" name="ooo_bundled_extensions">

<service
   name="application/ooo_bundled_extensions"
   type="service"
   version="1">

  <create_default_instance enabled="true" />
  <single_instance/>

  <dependency name="service_available"
          grouping="require_all"
          restart_on="none"
          type="service">
    <service_fmri value="svc:/system/filesystem/local:default" />
  </dependency>

  <exec_method
     type="method" name="start"
     exec="/lib/svc/method/svc-ooo_bundled_extensions start"
     timeout_seconds="160">
    <method_context>
      <method_credential user="root" group="root" />
    </method_context>
  </exec_method>


  <exec_method
     type="method" name="stop"
     exec="/lib/svc/method/svc-ooo_bundled_extensions stop"
     timeout_seconds="160" >
    <method_context>
      <method_credential user="root" group="root" />
    </method_context>
  </exec_method>

  <property_group name="startd" type="framework">
    <propval name="duration" type="astring"
         value="transient" />
    <propval name="ignore_error"
         type="astring" value="core,signal" />

  </property_group>

  <stability value="Unstable" />

  <template>
    <common_name>
      <loctext xml:lang="C">
    Installation of OpenOffice.org's bundled extensions.
      </loctext>
    </common_name>
  </template>
</service>

</service_bundle>
